# root/etc/sv/rslsync/run
#!/bin/sh
exec 2>&1
su -c "/usr/local/bin/rslsync --config /etc/rslsync/rslsync.conf --nodaemon" -s /bin/sh root